# Notes/snippets/launchd/com.huangrt.csnotes.ai-hotspots.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.huangrt.csnotes.ai-hotspots</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>/Users/bytedance/CS-Notes/Notes/snippets/generate-ai-hotspots-daily.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/bytedance/CS-Notes</string>

  <key>RunAtLoad</key>
  <true/>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>8</integer>
    <key>Minute</key>
    <integer>30</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/bytedance/CS-Notes/.local/ai-hotspots/logs/launchd.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/bytedance/CS-Notes/.local/ai-hotspots/logs/launchd.err.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>REPO_ROOT</key>
    <string>/Users/bytedance/CS-Notes</string>
    <key>AI_HOTSPOTS_NOTIFY</key>
    <string>1</string>
  </dict>
</dict>
</plist>
